# tools/qa/com.fmn.snapshot-diff-regression.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Unofficial FortiMonitor Toolkit - Gregori Jenkins
  FMN-222: launchd job that runs the snapshot diff regression spec once a day.

  Install:
    cp tools/qa/com.fmn.snapshot-diff-regression.plist \
       ~/Library/LaunchAgents/com.fmn.snapshot-diff-regression.plist
    launchctl load ~/Library/LaunchAgents/com.fmn.snapshot-diff-regression.plist

  Trigger immediately for a smoke test:
    launchctl start com.fmn.snapshot-diff-regression

  Disable / remove:
    launchctl unload ~/Library/LaunchAgents/com.fmn.snapshot-diff-regression.plist

  Constraints:
    - Requires `node tools/dev/launcher.mjs` to be running with the operator
      signed into FortiMonitor (the spec connects via CDP to localhost:9222).
      The runner will exit with a useful error if the launcher is not up.
    - Mac must be awake at the scheduled time. launchd will catch up after
      wake but only if a scheduled fire was missed; "skipped" is silent.
    - Logs land under tools/qa/logs/ (gitignored).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fmn.snapshot-diff-regression</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/gregorijenkins/Projects/unofficial-fortimonitor-browser-plugin/tools/qa/snapshot-diff-regression-runner.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/gregorijenkins/Projects/unofficial-fortimonitor-browser-plugin</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/gregorijenkins</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>7</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/gregorijenkins/Projects/unofficial-fortimonitor-browser-plugin/tools/qa/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/gregorijenkins/Projects/unofficial-fortimonitor-browser-plugin/tools/qa/logs/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
